# files/system/usr/lib/systemd/system/docker.socket.d/10-hyprblue.conf
[Unit]
After=hyprblue-groups.service
Requires=hyprblue-groups.service
